# services/templates/jarvis-self-healer.service
[Unit]
Description=Jarvis Self-Heal Daemon (monitor & auto-restart)
After=default.target

[Service]
Type=simple
ExecStart=__VENV_PY__ __JARVIS_ROOT__/services/self_healer.py
Restart=on-failure
RestartSec=15s
Nice=15
Environment="PYTHONPATH=__JARVIS_ROOT__"
Environment="JARVIS_ROOT=__JARVIS_ROOT__"
Environment="VAULT_ROOT=__VAULT_ROOT__"
WorkingDirectory=__JARVIS_ROOT__

[Install]
WantedBy=default.target
