# services/templates/jarvis-context-updater.service
[Unit]
Description=Jarvis Context Updater (run by timer)

[Service]
Type=oneshot
ExecStart=__VENV_PY__ __JARVIS_ROOT__/services/context_updater.py
Nice=15
Environment="PYTHONPATH=__JARVIS_ROOT__"
Environment="JARVIS_ROOT=__JARVIS_ROOT__"
Environment="VAULT_ROOT=__VAULT_ROOT__"
WorkingDirectory=__JARVIS_ROOT__
